# wifi-notifier.service.template
[Unit]
Description=WiFi Proximity Notifier
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=__INSTALL_DIR__
ExecStart=__INSTALL_DIR__/venv/bin/python3 app.py
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/__UID__/bus
Environment=XDG_RUNTIME_DIR=/run/user/__UID__
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
